# systemd/system/bazalt-launch.target
[Unit]
Description=Bazalt Launch Target
